# dd-edge/install/dd-edge.timer
[Unit]
Description=DD-Edge runs every 5 minutes

[Timer]
OnBootSec=30sec
OnUnitActiveSec=5min

[Install]
WantedBy=multi-user.target
